# config/snivirtualproxy.service
[Unit]
Description=SNI Virtual Proxy

[Service]
ExecStart=/usr/local/sbin/snivirtualproxy

[Install]
WantedBy=multi-user.target